# overlay/etc/init/init-gogs.conf
# init Gogs configuration

description "initializes Gogs on first start"
author "Scaleway <opensource@scaleway.com>"

start on (startup)

task

script
	/usr/local/sbin/oc-init-gogs
end script
